# files/customization/xdg-base-dir/usr/lib/systemd/user/gpg-agent-setup.timer
[Unit]
Description=Timer for gpg-agent-setup

[Timer]
OnStartupSec=10

[Install]
WantedBy=timers.target
