# files/testing/customization/xdg-base-dir/usr/lib/systemd/user/gpg-agent-setup.service
[Unit]
Description=Setup GnuPG directory
Wants=network-online.target
Before=gpg-agent.socket

[Service]
Type=oneshot
ExecStart=/usr/libexec/grand-os/gpg-agent-setup
