# files/testing/customization/xdg-base-dir/usr/lib/systemd/user/gpg-agent-setup.timer
[Unit]
Description=Timer for gpg-agent-setup

[Timer]
OnBootSec=10

[Install]
WantedBy=timers.target
