# files/system/common/usr/lib/systemd/user/flatpak-bundle-install.timer
[Unit]
Description=Install Flatpak bundles Trigger

[Timer]
OnBootSec=60sec

[Install]
WantedBy=timers.target
