# files/system/common/usr/lib/systemd/user/flatpak-override-update.timer
[Unit]
Description=Update Flatpak overrides Trigger

[Timer]
OnBootSec=2sec

[Install]
WantedBy=timers.target
