# files/system/common/usr/lib/systemd/user/flatpak-override-update.service
[Unit]
Description=Update Flatpak overrides

[Service]
Type=oneshot
ExecStart=/usr/libexec/flatpak-override-update
