# files/system/common/usr/lib/systemd/user/gnupg-unit-update.timer
[Unit]
Description=Update GnuPG unit files Trigger

[Timer]
OnBootSec=30sec

[Install]
WantedBy=timers.target
