# files/system/common/usr/lib/systemd/user/gnupg-unit-update.service
[Unit]
Description=Update GnuPG unit files

[Service]
Type=oneshot
ExecStart=/usr/libexec/grand-os/gnupg-unit-update
