# files/testing/security/kernel/usr/lib/systemd/user/flatpak-harden.timer
[Unit]
Description=Timer for flatpak-harden

[Timer]
OnBootSec=30

[Install]
WantedBy=timers.target
