# files/testing/customization/default-application/usr/lib/systemd/user/flatpak-bundle-setup.timer
[Unit]
Description=Timer for flatpak-bundle-setup

[Timer]
OnBootSec=30

[Install]
WantedBy=timers.target
