# files/testing/customization/xdg-base-dir/etc/skel/.config/systemd/user/gpg-agent-extra.socket.d/override.conf
[Socket]
ListenStream=
ListenStream=%t/gnupg/S.gpg-agent.extra
